# debian/hermes.prerm
#!/bin/bash

set -e

if [ "$1" = remove ]; then
    pam-auth-update --package --remove hermes
    systemctl stop hermes
    systemctl disable hermes
fi
