# debian/hermes.postinst
#!/bin/bash

set -e

pam-auth-update --package

addgroup hermes || true

systemctl enable hermes
systemctl start hermes
