# debian/hermes.postinst
#!/bin/bash

set -e

pam-auth-update --package

addgroup --system hermes

systemctl enable hermes
systemctl start hermes
